# systemd/aud_arb.service
[Unit]
Description=AUD Arbitrage Detector (Milestone 1 - detection only)
After=network.target

[Service]
User=mertcan
WorkingDirectory=/opt/aud_arb
Environment="PAIR=BTC/AUD"
Environment="CSV_PATH=/opt/aud_arb/out/ir_boa.csv"
ExecStart=/opt/aud_arb/.venv/bin/python /opt/aud_arb/src/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
